# src/daumalapps.service
[Unit]
After=network.service

[Service]
User=root
ExecStart=/usr/local/bin/daumalapps.sh

[Install]
WantedBy=default.target